# DEBIAN/postinst
#!/bin/sh
name='remote-dconf'
chown -R root:root /opt/$name
systemctl daemon-reload
systemctl enable $name || true
systemctl restart $name || true
