# DEBIAN/postinst
#!/bin/sh
name='basic-file-server'
chown -R root:root /opt/$name
systemctl daemon-reload
systemctl enable $name || true
systemctl restart $name || true
